# challenges/misc/mission-impossible/source/server/supervisord-handler.ini
[program:handler]
directory=/app
command=/usr/local/bin/python3 /app/handler.py
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs = 0
autorestart=false
